# limonada.cron
#!/bin/bash

find /path/limonada/media/tmp/ -mmin +60 -mindepth 1 -exec rm -rf {} \;


