# limonada.cron
#!/bin/bash

find /path/limonada/media/tmp/ -mmin +60 -exec rm -rf {} \;


